# scripts/twinshare-rest-api.service
[Unit]
Description=twinshare REST API Server
After=network.target

[Service]
Type=forking
User=$USER
Group=$USER
WorkingDirectory=/home/tom/github/twinshare/share
Environment="PYTHONPATH=/home/tom/github/twinshare/share"
ExecStartPre=/bin/mkdir -p /var/log/twinshare /var/run/twinshare
ExecStart=/home/tom/github/twinshare/share/venv/bin/python3 /home/tom/github/twinshare/share/scripts/start_rest_api.py start --log-file /var/log/twinshare/rest_api.log --pid-file /var/run/twinshare/rest_api.pid
ExecStop=/home/tom/github/twinshare/share/venv/bin/python3 /home/tom/github/twinshare/share/scripts/start_rest_api.py stop --pid-file /var/run/twinshare/rest_api.pid
PIDFile=/var/run/twinshare/rest_api.pid
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
